# installers/debian/keychainbackup-0.3/DEBIAN/postinst
#!/bin/bash
chmod 644 /etc/systemd/system/keychainmanager.service
systemctl daemon-reload
systemctl enable keychainmanager.service
systemctl start keychainmanager
